# tools/runner_provision/supervisord-gcp.conf.tmpl
[supervisord]
user=root

[program:runner]
command=/startup_script.py $labels